# infra/systemd/aether-boop-30min.service
[Unit]
Description=Aether Autonomous BOOP Group — 30min cadence (engineering-flow-check, delegation-enforcer, telegram-health-boop)
Documentation=file:///home/jared/projects/AI-CIV/aether/tools/boop_group_runner.sh
Documentation=file:///home/jared/projects/AI-CIV/aether/.claude/memory/agent-learnings/devops-engineer/2026-02-23--autonomous-boop-scheduling-system.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jared
WorkingDirectory=/home/jared/projects/AI-CIV/aether
Environment="PATH=/home/jared/.local/bin:/home/jared/.nvm/versions/node/v22.14.0/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/jared"
ExecStart=/home/jared/projects/AI-CIV/aether/tools/boop_group_runner.sh 30min
StandardOutput=append:/home/jared/projects/AI-CIV/aether/logs/boops/systemd-30min.log
StandardError=append:/home/jared/projects/AI-CIV/aether/logs/boops/systemd-30min.log
TimeoutStartSec=20min

# Hardening (mirrors aether-overnight-monitor.service style)
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/jared/projects/AI-CIV/aether/logs /home/jared/projects/AI-CIV/aether/inbox /home/jared/projects/AI-CIV/aether/.claude /home/jared/projects/AI-CIV/aether/memories /home/jared/projects/AI-CIV/aether/.current_session
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
