# infra/systemd/aether-boop-30min.timer
[Unit]
Description=Aether BOOP 30min timer — fires aether-boop-30min.service on the :00 and :30 of every hour
Documentation=file:///home/jared/projects/AI-CIV/aether/.claude/memory/agent-learnings/devops-engineer/2026-02-23--autonomous-boop-scheduling-system.md

[Timer]
OnCalendar=*:0/30
AccuracySec=30s
Persistent=true
Unit=aether-boop-30min.service

[Install]
WantedBy=timers.target
