# config/aether-knowledge-sync.service
[Unit]
Description=Aether Knowledge Base Sync
After=network.target

[Service]
Type=oneshot
User=jared
WorkingDirectory=/home/jared/projects/AI-CIV/aether
ExecStart=/usr/bin/python3 /home/jared/projects/AI-CIV/aether/tools/sync_knowledge.py
StandardOutput=append:/home/jared/projects/AI-CIV/aether/logs/knowledge_sync.log
StandardError=append:/home/jared/projects/AI-CIV/aether/logs/knowledge_sync.log

[Install]
WantedBy=multi-user.target
