# config/aether-session.service
[Unit]
Description=Aether Session Manager - keeps Claude Code + Telegram bridge always running
After=network.target

[Service]
Type=simple
User=jared
WorkingDirectory=/home/jared/projects/AI-CIV/aether
ExecStart=/home/jared/projects/AI-CIV/aether/tools/aether_session_manager.sh
Restart=always
RestartSec=10
StandardOutput=append:/home/jared/projects/AI-CIV/aether/logs/session_manager.log
StandardError=append:/home/jared/projects/AI-CIV/aether/logs/session_manager.log

[Install]
WantedBy=multi-user.target
