# config/aether-comms-gateway.service
[Unit]
Description=Pure Technology Communications Gateway
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=jared
WorkingDirectory=/home/jared/projects/AI-CIV/aether/tools/comms-gateway
ExecStart=/usr/bin/python3 /home/jared/projects/AI-CIV/aether/tools/comms-gateway/main.py
Restart=always
RestartSec=5
StandardOutput=append:/home/jared/projects/AI-CIV/aether/logs/comms_gateway.log
StandardError=append:/home/jared/projects/AI-CIV/aether/logs/comms_gateway.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
