# config/agentmail-monitor.service
[Unit]
Description=Aether AgentMail Monitor - polls aether-aiciv@agentmail.to every 30s
After=network.target
Wants=network.target

[Service]
Type=simple
User=jared
WorkingDirectory=/home/jared/projects/AI-CIV/aether
ExecStart=/usr/bin/python3 /home/jared/projects/AI-CIV/aether/tools/agentmail_monitor.py
Restart=always
RestartSec=15
StandardOutput=append:/home/jared/projects/AI-CIV/aether/logs/agentmail_monitor.log
StandardError=append:/home/jared/projects/AI-CIV/aether/logs/agentmail_monitor.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
