# config/blog-distribution.service
[Unit]
Description=Blog Distribution Pipeline - Check for new posts
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/jared/projects/AI-CIV/aether
ExecStart=/home/jared/projects/AI-CIV/aether/venv/bin/python tools/blog_distribution_pipeline.py check
User=jared
Environment="PATH=/home/jared/projects/AI-CIV/aether/venv/bin:/usr/bin"

[Install]
WantedBy=multi-user.target
