# config/gdrive-tester-monitor.service
[Unit]
Description=Aether Google Drive Tester Feedback Monitor
After=network.target

[Service]
Type=oneshot
User=jared
WorkingDirectory=/home/jared/projects/AI-CIV/aether
ExecStart=/usr/bin/python3 /home/jared/projects/AI-CIV/aether/tools/gdrive_tester_monitor.py check
StandardOutput=append:/home/jared/projects/AI-CIV/aether/logs/gdrive_tester_monitor.log
StandardError=append:/home/jared/projects/AI-CIV/aether/logs/gdrive_tester_monitor.log

[Install]
WantedBy=multi-user.target
