# config/paypal-webhook.service
[Unit]
Description=PayPal Auto-Split Webhook Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jared
Group=jared
WorkingDirectory=/home/jared/projects/AI-CIV/aether
ExecStart=/usr/bin/python3 /home/jared/projects/AI-CIV/aether/tools/paypal_auto_split.py --webhook
Restart=always
RestartSec=10
StandardOutput=append:/home/jared/projects/AI-CIV/aether/logs/paypal_webhook_service.log
StandardError=append:/home/jared/projects/AI-CIV/aether/logs/paypal_webhook_service.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
