# tools/gmail_monitor.service
[Unit]
Description=Aether Gmail Monitor Daemon
After=network.target

[Service]
Type=simple
User=jared
WorkingDirectory=/home/jared/projects/AI-CIV/aether
ExecStart=/usr/bin/python3 /home/jared/projects/AI-CIV/aether/tools/gmail_monitor.py daemon --interval 5
Restart=always
RestartSec=60
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=append:/home/jared/projects/AI-CIV/aether/logs/gmail_monitor_service.log
StandardError=append:/home/jared/projects/AI-CIV/aether/logs/gmail_monitor_service.log

[Install]
WantedBy=multi-user.target
